# overlay/etc/supervisor/conf.d/aria2.conf
[program:aria2]
command=/usr/bin/aria2c --conf-path=/home/aria/aria2.conf
user=aria
autostart=true
autorestart=true
startsecs=3
startretries=3
stdout_logfile=/dev/null